# etc/crontabs/root
45 3 * * * /etc/init.d/njit restart
*/20 * * * * cat /dev/null > /tmp/njit.log

