# etc/init.d/ethtool
#!/bin/sh /etc/rc.common
START=18
start() {
	echo "Setting eth0 to 100/full"
	mii-tool -A 100baseTx-FD eth0;
}
